# dashboard.ini
[program:httpd]
directory=/
command=/usr/sbin/httpd -DFOREGROUND
user=apache
group=apache
autostart=true
autorestart=true
priority=999
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log