# dashboard.ini
[program:apache2]
directory=/
command=/usr/sbin/apachectl -k start -DFOREGROUND
user=root
autostart=true
autorestart=true
priority=999
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log